# reports_auto/pack_stage/20250916T051048/reports_auto/export/20250916T005124/stage/tools/schedule/sma-oneclick.timer
[Unit]
Description=Run Smart Mail Agent One-Click Daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
Unit=sma-oneclick.service

[Install]
WantedBy=timers.target
