# reports_auto/pack_stage/20250916T051048/reports_auto/export/20250916T014311/bundle/tools/schedule/sma-oneclick.service
[Unit]
Description=Smart Mail Agent One-Click Run

[Service]
Type=oneshot
WorkingDirectory=%h/projects/smart-mail-agent_ssot
Environment=PYTHONNOUSERSITE=1
Environment=PYTHONPATH=src
Environment=SMA_SMTP_MODE=outbox
Environment=SMA_LLM_PROVIDER=none
ExecStart=%h/projects/smart-mail-agent_ssot/tools/one_click_s2s3s5.sh
